# etc/supervisor/conf.d/plasmashell.conf
[program:plasmashell]
command=/usr/bin/plasmashell --replace
autostart=true
autorestart=true
stderr_logfile=/var/log/plasmashell.err.log
stdout_logfile=/var/log/plasmashell.out.log
