# etc/init.d/mountbinds
#!/sbin/openrc-run

description="Asthethic FHS for Nitrux OS"

check_casper_script() {
    if [ -f /usr/share/initramfs-tools/scripts/casper-bottom/12fstab ]; then
        return 0 # File exists
    else
        return 1 # File does not exist
    fi
}

start() {
    ebegin "Checking prerequisites"
    if check_casper_script; then
        einfo "Casper fstab script detected, skipping mount binds; this is a Live session."
        return 0
    fi
    einfo "No Casper fstab script detected, proceeding with mount binds; this isn't a Live session."

    ebegin "Binding FHS directories"

    mount --bind /dev /System/Devices
    mount --bind /boot /System/Boot
    mount --bind /etc /System/Configuration
    mount --bind /usr/sbin /System/Binaries/Administrative
    mount --bind /usr/bin /System/Binaries/Standard
    mount --bind /opt /System/Binaries/Optional
    mount --bind /mnt /System/Filesystems/Internal
    mount --bind /media /System/Filesystems/External
    mount --bind /usr/lib /System/Libraries
    mount --bind /usr/share /System/Shareable
    mount --bind /var /System/Variable
    mount --bind /home /Users

    eend $?
}

stop() {
    ebegin "Unbinding FHS directories"

    umount /System/Devices
    umount /System/Boot
    umount /System/Configuration
    umount /System/Binaries/Administrative
    umount /System/Binaries/Standard
    umount /System/Binaries/Optional
    umount /System/Filesystems/External
    umount /System/Filesystems/Internal
    umount /System/Libraries
    umount /System/Shareable
    umount /System/Variable
    umount /Users

    eend $?
}
