# debian/postinst
#!/bin/bash

echo QT_STYLE_OVERRIDE=kvantum >> /etc/environment 

rm -R /usr/share/icons/breeze_cursors /usr/share/icons/Breeze_Snow
rm /usr/share/applications/calamares.desktop /usr/share/applications/org.kde.kdeconnect.sms.desktop /usr/share/applications/org.kde.kdeconnect_open.desktop /usr/share/applications/org.kde.kdeconnect.app.desktop /usr/share/applications/htop.desktop /usr/share/applications/mc.desktop /usr/share/applications/mcedit.desktop /usr/share/applications/nsnake.desktop /usr/share/applications/org.kde.kinfocenter.desktop
ln -svf /usr/games/nsnake /bin/nsnake

git clone https://github.com/robbyrussell/oh-my-zsh.git /etc/skel/.oh-my-zsh

rm /bin/sh.distrib
# /usr/bin/dpkg --remove --no-triggers --force-remove-essential --force-bad-path dash &> /dev/null
ln -sv /bin/mksh /bin/sh
# /usr/bin/dpkg --remove --no-triggers --force-remove-essential --force-bad-path dash &> /dev/null

sed -i 's+SHELL=/bin/sh+SHELL=/bin/zsh+g' /etc/default/useradd
sed -i 's+DSHELL=/bin/bash+DSHELL=/bin/zsh+g' /etc/adduser.conf

sed -i 's/#DefaultTimeoutStartSec=90s/DefaultTimeoutStartSec=5s/g' /etc/systemd/system.conf
sed -i 's/#DefaultTimeoutStopSec=90s/DefaultTimeoutStopSec=5s/g' /etc/systemd/system.conf

mkdir -p /etc/systemd/system/networking.service.d/
bash -c 'echo -e "[Service]\nTimeoutStartSec=20sec" > /etc/systemd/system/networking.service.d/timeout.conf'

systemctl mask avahi-daemon.service
systemctl disable cupsd.service cupsd-browsed.service NetworkManager-wait-online.service keyboard-setup.service

sed -i 's/ACTION!="add", GOTO="libmtp_rules_end"/ACTION!="bind", ACTION!="add", GOTO="libmtp_rules_end"/g' /lib/udev/rules.d/69-libmtp.rules

mkdir -p /Applications
mkdir -p /etc/skel/Applications
mkdir -p /etc/skel/.local/bin

echo -e "PATH=$PATH:/Applications\n" > /etc/environment
sed -i "s|secure_path\=.*$|secure_path=\"$PATH:/Applications\"|g" /etc/sudoers
sed -i "/env_reset/d" /etc/sudoers

APPS_SYS='
https://github.com/AppImage/appimaged/releases/download/continuous/appimaged-x86_64.AppImage
https://raw.githubusercontent.com/Nitrux/storage/master/AppImages/appimage-cli-tool-x86_64.AppImage
'

APPS_USR='
'

for x in $APPS_SYS; do
    axel -a -q -o /Applications $x
done

for x in $APPS_USR; do
    axel -a -q -o /Applications $x
done

chmod +x /Applications/*

mv /Applications/appimage-cli-tool-x86_64.AppImage /Applications/app

mv /Applications/appimaged-x86_64.AppImage /etc/skel/.local/bin/appimaged

ln -svf /usr/share/kwin/scripts/krohnkite/metadata.desktop /usr/share/kservices5/krohnkite.desktop

git clone https://github.com/powerline/fonts /tmp/powerline-fonts
cp -r /tmp/powerline-fonts/Terminus/PSF/*.psf.gz /usr/share/consolefonts
rm -r /tmp/powerline-fonts

rm /etc/default/console-setup && touch /etc/default/console-setup

>> /etc/default/console-setup printf "%s\n" \
	"+# CONFIGURATION FILE FOR SETUPCON" \
	" " \
	"# Consult the console-setup(5) manual page." \
	" " \
	'ACTIVE_CONSOLES="/dev/tty[1-6]"' \
	'CHARMAP="UTF-8"' \
	" " \
	"VIDEOMODE=" \
	'FONT="ter-powerline-v16n.psf.gz"' \

rm /etc/skel/.config/Trolltech.conf && touch /etc/skel/.config/Trolltech.conf

>> /etc/skel/.config/Trolltech.conf  printf "%s\n" \
	"[qt]" \
	"KDE\contrast=0" \
	"KWinPalette\activeBackground=#2b2c31" \
	"KWinPalette\activeBlend=#f9f9f9" \
	"KWinPalette\activeForeground=#f9f9f9" \
	"KWinPalette\activeTitleBtnBg=#2b2c31" \
	"KWinPalette\frame=#2b2c31" \
	"KWinPalette\inactiveBackground=#37383f" \
	"KWinPalette\inactiveBlend=#4b4743" \
	"KWinPalette\inactiveForeground=#9bacbb" \
	"KWinPalette\inactiveFrame=#2c2d32" \
	"KWinPalette\inactiveTitleBtnBg=#2c2d32" \
	"Palette\active=#ffffff, #505659, #37383f, #313238, #1e1e22, #27282c, #ffffff, #ffffff, #ffffff, #3a3f41, #2b2c31, #151618, #26c6da, #eff0f1, #2980b9, #7f8c8d, #eff0f1, #000000, #31363b, #eff0f1, #ffffff" \
	"Palette\disabled=#ffffff, #505659, #37383f, #313238, #1e1e22, #27282c, #ffffff, #ffffff, #ffffff, #3a3f41, #2b2c31, #151618, #2b2c31, #ffffff, #2980b9, #7f8c8d, #eff0f1, #000000, #31363b, #eff0f1, #ffffff" \
	"Palette\inactive=#fafafa, #515659, #383a40, #333439, #1f1f23, #28292d, #fafafa, #ffffff, #fbfbfb, #3b4042, #2c2d32, #161719, #28c4d8, #def0f3, #377faf, #7f8b8c, #ecedee, #000000, #31363b, #eff0f1, #ffffff" \
