# etc/skel/.config/autostart/nx-welcome-wizard.desktop
[Desktop Entry]
Type=Application
Name=Welcome Wizard
Exec=nx-welcome-wizard
Icon=start-here
Comment=
