# etc/xdg/autostart/appimaged.desktop
[Desktop Entry]
Type=Application
Name=AppImage Desktop Integration Daemon
Exec=$HOME/.local/bin/appimaged
Icon=start-here
Comment=
