# etc/xdg/autostart/nomad-welcome-wizard.desktop
[Desktop Entry]
Type=Application
Name=Welcome Wizard
Exec=nomad-welcome-wizard
Icon=
Comment=<optinal comments>
